# bin/zenserver.service
[Unit]
Description=Server for the ZenLighting project

[Service]
Type=simple
ExecStart=/usr/bin/env python3 /insert/path/to

[Install]
WantedBy=multi-user.target